# cookbooks/ecom-platform/templates/default/dj/dj.erb
#!/bin/bash

### BEGIN INIT INFO
# Provides:  dj
# Required-Start:$all
# Required-Stop: $all
# Default-Start: 2 3 4 5
# Default-Stop:  0 1 6
# Short-Description: dj application service
# Description:  dj application service
### END INIT INFO

APP_USER="<%=@user%>"
APP_GROUP="<%=@group%>"
APP="dj"
APP_PATH="<%=@project_path%>/current"

eecute() {
  sudo -u $APP_USER -H bash -l -c "$1"
  }

start() {
    cd $APP_PATH
    su - $APP_USER -c "cd $APP_PATH;source /etc/default/<%=@app%>.conf;RAILS_ENV=production bundle exec bin/delayed_job --queues=<%=@queues%> -n <%=@count%> start"

    if [ $? -eq 0 ]; then
        sleep 2
        echo "$APP started successfully with PID "
    else
        echo "$APP could not start. Please check logs."
        fi
}

stop() {
      if [ $(ps aux |grep 'delayed' | grep -v grep | wc -l) -eq 0 ];then
          echo "Application $APP is already stopped ";
      else
         echo "Stopping $APP process"
         su - $APP_USER -c "cd $APP_PATH;source /etc/default/<%=@app%>.conf;RAILS_ENV=production bundle exec bin/delayed_job stop"
      fi
}

case "$1" in
  start)
    echo "Starting $APP"
    start
    echo "$APP started."
    ;;

  stop)
    echo "Stopping $APP"
    stop
    echo "$APP stopped."
    ;;

  restart)
    echo  "Restarting $APP."
    stop
 #   sleep 2
    start
    echo "$APP restarted."
    ;;

  status)
    if [ $(ps aux |grep 'delayed' | grep -v grep | wc -l) -eq 0 ];then
      echo  "Application $APP is stopped."
    else
      echo  "Application $APP is running with pid: " `cat $PID_FILE`
    fi
    ;;

  *)
    service_name=/etc/init.d/$APP
    echo "Usage: $service_name {start|stop|restart|status}" >&2
    exit 1
    ;;
esac

exit 0
